# dot_config/systemd/user/aw-awatcher.service
[Unit]
Description=awatcher — window + AFK watcher (Hyprland)
# NOT wired into default.target — starting on boot causes a race with
# Hyprland's initial toplevel batch, after which awatcher latches onto
# stale zwlr_foreign_toplevel_handle_v1 handles and wedges into a
# permanent "Current window is unknown" loop. Process stays alive so
# Restart=always never fires.
#
# Started instead from Hyprland exec-once (see hyprland.conf), AFTER
# `dbus-update-activation-environment` has imported WAYLAND_DISPLAY and
# HYPRLAND_INSTANCE_SIGNATURE into the systemd user environment. By that
# point Hyprland has already announced its initial toplevel batch and
# awatcher latches onto live handles.
After=aw-server.service
Requires=aw-server.service

[Service]
Type=simple
ExecStart=/usr/bin/aw-awatcher --idle-timeout 300
Restart=always
RestartSec=5
