# dot_config/systemd/user/aw-awatcher.service
[Unit]
Description=awatcher — window + AFK watcher (Hyprland)
After=aw-server.service graphical-session.target
Requires=aw-server.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/aw-awatcher --idle-timeout 300
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
