# dot_config/systemd/user/aw-health-check.service
[Unit]
Description=ActivityWatch stack health check
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/dev/toggl-tracker/alerting/aw-health-check.sh

[Install]
WantedBy=default.target
