# dot_config/systemd/user/aw-today-cache.timer
[Unit]
Description=Refresh aw-today JSON cache every 10 minutes
After=graphical-session.target
PartOf=graphical-session.target

[Timer]
OnActiveSec=30s
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
