# dot_config/systemd/user/aw-today-cache.timer
[Unit]
Description=Refresh aw-today JSON cache every 3 minutes
After=graphical-session.target
PartOf=graphical-session.target

[Timer]
OnActiveSec=30s
OnUnitActiveSec=3min
AccuracySec=15s
Persistent=false

[Install]
WantedBy=timers.target
